# __dht.conf
[program:dht]
command=/usr/bin/python /usr/local/sbin/rpimonitor/dht.py 2302 21
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/dht.log
stderr_logfile=/var/log/dht.log
